# app/tools/electron-checkBalance.conf
#KHD: June 08, 2017
#Instruction: copy this script to /etc/init/ directory
description "Start the Electron Checkbalance app in fullscreen mode"
author "Khoa Do <khoa.do@nkidcorp.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process
script
  cd /home/pi/devevopment/checkbalance-app/
  exec sudo npm start
end script